# src/.config/systemd/user/delete-tmp-media.timer
[Unit]
Description=Empty media folders monthly

[Timer]
OnCalendar=monthly
Persistent=True

[Install]
WantedBy=timers.target
